# php-xingstation-api/supervisor/lumen-worker-common.ini
[program:common]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --queue=warning --tries=3
autostart=true
autorestart=true
numprocs=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue.log